# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/tinai/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── llama-server ──────────────────────────────────────────────────
[program:llama]
command=/usr/local/bin/start-llama.sh
autostart=true
autorestart=true
priority=10
startretries=3
startsecs=5
stdout_logfile=/var/log/tinai/llama.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/tinai/llama.err
stderr_logfile_maxbytes=5MB
environment=HOME="/root",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# ── Open WebUI ────────────────────────────────────────────────────
[program:webui]
command=/usr/local/bin/start-webui.sh
autostart=true
autorestart=true
priority=20
startretries=3
startsecs=10
stdout_logfile=/var/log/tinai/webui.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/tinai/webui.err
stderr_logfile_maxbytes=5MB
environment=HOME="/root",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# ── OpenFang (binaire Rust natif) ────────────────────────────────
[program:openfang]
command=/usr/local/bin/start-openfang.sh
autostart=true
autorestart=true
priority=30
startretries=3
startsecs=15
stdout_logfile=/var/log/tinai/openfang.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/tinai/openfang.err
stderr_logfile_maxbytes=5MB
environment=HOME="/root",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
